# build_files/files/usr/lib/systemd/system/install-default-flatpaks.timer
[Unit]
Description=Run default Flatpak installer shortly after boot

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target